# crontab/user
*/10 * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) /usr/bin/subtube update --secret
*/30 * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) notify-send -t 1500 "$(date +\%H:\%M)"
* * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) /home/infiniter/.scripts/torrent_notifier
